# com.user.obsidian.backup.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
    <dict>
        <key>Label</key>
        <string>com.user.obsidian.backup</string>
        <key>ProgramArguments</key>
        <array>
            <string>/bin/bash</string>
            <string>/Users/TODO_YOURUSERNAME/obsidian/auto-backup-obsidian.sh</string>
        </array>
        <key>StartInterval</key>
        <integer>1800</integer>
        <key>RunAtLoad</key>
        <true/>
        <key>StandardErrorPath</key>
        <string>/Users/TODO_YOURUSERNAME/obsidian/backup-error.log</string>
        <key>StandardOutPath</key>
        <string>/Users/TODO_YOURUSERNAME/obsidian/backup-output.log</string>
        <key>EnvironmentVariables</key>
        <dict>
            <key>PATH</key>
            <string>/usr/local/bin:/usr/bin:/bin:/usr/sbin:/sbin</string>
        </dict>
    </dict>
</plist>